# deploy/com.agentradar.daily.plist
<?xml version="1.0" encoding="UTF-8"?>
<!-- TEMPLATE — `scripts/install-launchd.sh` replaces __REPO__ with this repo's absolute path
     and installs the result to ~/Library/LaunchAgents (never committed). Daily at 08:30. -->
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
  <key>Label</key><string>com.agentradar.daily</string>
  <key>ProgramArguments</key>
  <array>
    <string>/bin/bash</string>
    <string>__REPO__/scripts/run-daily.sh</string>
  </array>
  <key>WorkingDirectory</key><string>__REPO__</string>
  <key>StartCalendarInterval</key>
  <dict>
    <key>Hour</key><integer>8</integer>
    <key>Minute</key><integer>30</integer>
  </dict>
  <key>RunAtLoad</key><false/>
  <key>StandardOutPath</key><string>__REPO__/data/state/launchd-daily.log</string>
  <key>StandardErrorPath</key><string>__REPO__/data/state/launchd-daily.log</string>
</dict>
</plist>
